# user.service
# SPDX-License-Identifier: LGPL-2.1-or-later
#
# This file is part of systemd-user-sessions.
#
# systemd is free software; you can redistribute it and/or modify it
# under the terms of the GNU Lesser General Public License as published by
# the Free Software Foundation; either version 2.1 of the License, or
# (at your option) any later version.
#
# systemd-user-sessions — Manage user login sessions

[Unit]
Description=User Login Session Manager
Documentation=man:systemd-user-sessions.service(8)
After=network.target

[Service]
Type=simple
ExecStart=/usr/lib/systemd/systemd-user-service
StandardOutput=inherit
StandardError=journal
Restart=always
RestartSec=5
# Ensure /run/systemd/ exists for the Unix domain socket
RuntimeDirectory=systemd
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
